# SOURCES/dts-n62-issue-msg.timer
[Unit]
Description=Update the DTS N62 motd every 5 minutes

[Timer]
OnBootSec=0min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
